# update-motd.d/20-sysinfo
#!/bin/bash

# Load color palette file
source /etc/update-motd.d/theme

# OUTPUT
